# ha_logi_host/rootfs/etc/services.d/ha_logi_host/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Logi Host Switch..."

# Debug: show hidraw device permissions
bashio::log.info "hidraw devices:"
for dev in /dev/hidraw*; do
    if [ -e "$dev" ]; then
        bashio::log.info "  $(ls -la "$dev")"
    fi
done
bashio::log.info "Running as: $(id)"

export MQTT_HOST="$(bashio::services mqtt 'host')"
export MQTT_PORT="$(bashio::services mqtt 'port')"
export MQTT_USER="$(bashio::services mqtt 'username')"
export MQTT_PASS="$(bashio::services mqtt 'password')"

export RECEIVER_TYPE="$(bashio::config 'receiver_type')"
export LOG_LEVEL="$(bashio::config 'log_level')"

exec python3 -m logi_host
